# init.d/skynet-gentoo
#!/sbin/runscript

depend() {
        need local
}


PROGRAM="/etc/skynet-client/skynet-client-bot.pl"
start() {
        killall -9 skynet-client-bot.pl 2> /dev/null  || true
        ebegin "Starting Skynet client bot"
        start-stop-daemon --start -b --quiet --exec $PROGRAM
        eend $?
}
stop() {
        ebegin "Stopping Skynet client bot"
        killall -9 skynet-client-bot.pl 2> /dev/null
        eend $?

}

